# .config/systemd/user/paction.service
[Unit]
Description=Execute actions relative to process existance

[Service]
Type=simple
Restart=always
RestartSec=1
ExecStart=%h/.cargo/bin/paction

[Install]
WantedBy=default.target
